# image-spider/supervisord.conf
[program:imagespider]
command = /home/dotcloud/current/image-spider/main.py
