# spider/supervisord.conf
[program:imagespider]
command = python3.2 /home/dotcloud/current/main.py
